# dind.conf
[program:dind]
command=/usr/local/bin/wrapdocker
