# jenkins.conf
[program:jenkins]
command=/usr/local/bin/runjenkins
